# ops/dhanam-newwealth.service
[Unit]
Description=Dhanam New Wealth API
After=network.target postgresql.service

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/new-wealth/api
ExecStart=/usr/bin/node dist/index.js
Restart=always
RestartSec=3
MemoryMax=512M
Environment=NODE_ENV=production
Environment=FILE_STORAGE_DIR=/var/lib/dhanam-newwealth
EnvironmentFile=/home/ubuntu/new-wealth/api/.env
# systemd hardening — the app only writes to its upload dir + tmp.
ReadWritePaths=/var/lib/dhanam-newwealth /tmp
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only

[Install]
WantedBy=multi-user.target
